# src/etc/systemd/system/dnsmasq.service.d/gygeslink-ordering.conf
# GygesLink — Drop-in systemd pour dnsmasq
#
# Désactive le démarrage automatique de dnsmasq au boot.
# dnsmasq est lancé manuellement par gygeslink-network-setup.sh
# (en arrière-plan, pas via systemctl) pour éviter un deadlock :
# le drop-in After=gygeslink-network-setup empêchait dnsmasq de
# démarrer tant que network-setup tourne, mais network-setup
# appelait systemctl restart dnsmasq qui attendait la fin de
# network-setup = boucle infinie.
#
# Solution : dnsmasq est lancé directement par le script réseau,
# et ce drop-in empêche systemd de le démarrer au boot.

[Unit]
# Pas de dépendance sur gygeslink-network-setup — dnsmasq est
# géré manuellement par le script réseau.

[Service]
# Ne pas démarrer automatiquement au boot
# (le script réseau le lance quand usb0 est prêt)