# src/etc/systemd/system/dnsmasq.service.d/gygeslink-ordering.conf
# GygesLink — Drop-in systemd pour dnsmasq
#
# Garantit que dnsmasq démarre APRÈS que usb0 soit configuré
# (192.168.100.1/24 doit exister pour que dnsmasq puisse y écouter).
# Sans ce drop-in, dnsmasq pourrait démarrer avant usb0 et ignorer
# l'interface — le PC ne recevrait pas d'IP via DHCP.

[Unit]
After=gygeslink-network-setup.service
Wants=gygeslink-network-setup.service
