# src/etc/systemd/system/gygeslink-api.service
[Unit]
Description=GygesLink API REST pour l'app desktop
After=gygeslink-network-setup.service
Wants=gygeslink-network-setup.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/gygeslink-api.py
User=root
Restart=on-failure
RestartSec=5
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target