# src/etc/systemd/system/gygeslink-tor.service
[Unit]
Description=GygesLink — Tor + obfs4 (anonymisation transparente)
# Tor démarre après :
#   - iptables-drop (trafic bloqué pendant le démarrage de Tor)
#   - wireguard (tier status connu → torrc configuré correctement)
After=gygeslink-network-setup.service
After=gygeslink-wireguard.service

# Requires : si gygeslink-network-setup échoue, Tor ne démarre pas
# (pas de réseau = pas de Tor = trafic resté bloqué = fail-close correct)
Requires=gygeslink-network-setup.service

# Note : bridges.conf est optionnel — le script gygeslink-tor-prestart.sh
# active ou désactive UseBridges selon le tier et la présence des bridges.
# Tor peut démarrer sans bridges (Tier 1 direct, sans obfs4).

[Service]
# forking : Tor se daémonise lui-même (fork en arrière-plan après init).
# systemd attend que le processus parent se termine, ce qui indique
# que Tor a fini son initialisation.
Type=forking

# Tor tourne sous l'utilisateur dédié debian-tor (créé par le paquet tor)
# Cet utilisateur a uniquement les permissions nécessaires pour Tor.
User=debian-tor
Group=debian-tor

# Configurer OutboundBindInterface selon le tier (Tier 1 ou Tier 2)
# avant de lancer Tor. Ce script lit /run/gygeslink-tier.status.
# Préfixe + : ExecStartPre tourne en root (nécessaire pour sed -i sur torrc).
ExecStartPre=+/usr/local/bin/gygeslink-tor-prestart.sh

ExecStart=/usr/bin/tor -f /etc/tor/torrc

# PIDFile : indique à systemd où Tor écrit son PID (pour le suivi)
PIDFile=/run/tor/tor.pid

# Timeout généreux : les bridges obfs4 peuvent prendre du temps à répondre
# (surtout si le FAI pratique un DPI agressif).
# 90s était trop court (piège connu) → 180s
TimeoutStartSec=180

# Redémarrage en cas d'échec (ex: bridge temporairement indisponible)
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
